# pkgstats-cn.timer
[Unit]
Description=Report timer for pkgstats-cn

[Timer]
OnCalendar=weekly
Persistent=true
